# etc/services.d/avahi-daemon/run
#!/usr/bin/with-contenv bash

avahi-daemon
