# etc/services.d/smbd/run
#!/usr/bin/with-contenv bash

nmbd -D
smbd -FSs /data/smb.conf
